# templates/etc/systemd/system/roothints.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=Update root hints for unbound
After=network.target

[Service]
ExecStart=/usr/bin/curl -o {{ unbound_root_hints_dir }} {{ unbound_root_hints_url }}
